# service/filebrowser.service
[Unit]
Description=filebrowser
After=network.target

[Service]
ExecStart=filebrowser -c /etc/filebrowser.conf

[Install]
WantedBy=multi-user.target